# Caddy 2/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Caddy 2
# Runs Caddy
# ==============================================================================

bashio::log.info 'Starting Caddy...'

# Run Caddy
exec /usr/bin/caddy.sh